# dist/launchd/com.marshall.copad.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marshall.copad.daemon</string>

    <!--
        ProgramArguments runs copadd through a small wrapper
        (install-macos.sh installs it next to this plist's target) that
        sources ~/.config/copad/secrets.env before exec'ing
        ~/.cargo/bin/copadd. That file (mode 600, NOT overwritten by
        reinstall) carries plugin bearer tokens like
        COPAD_WEB_BRIDGE_TOKEN so they never live in this world-readable
        plist. The wrapper itself uses $HOME (launchd sets it for
        gui/<uid> agents); this path can't, so install-macos.sh rewrites
        HOME_PLACEHOLDER at copy time (launchd does not expand ~).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>HOME_PLACEHOLDER/.config/copad/copadd-launch.sh</string>
    </array>

    <!--
        RunAtLoad fires at launchd bootstrap (i.e. login). KeepAlive
        with no submap means "always running" — crash, manual kill,
        whatever; launchd restarts. Matches the systemd Restart=
        on-failure semantic close enough for parity.
    -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!--
        Pipe stdout + stderr to the user log dir so `tail
        ~/Library/Logs/copad-daemon.{out,err}.log` is the macOS
        equivalent of journalctl user-level lookup for copad-daemon.
    -->
    <key>StandardOutPath</key>
    <string>HOME_PLACEHOLDER/Library/Logs/copad-daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_PLACEHOLDER/Library/Logs/copad-daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <!--
        ProcessType=Interactive keeps the daemon at the same nice
        level as Aqua apps. Background would deprioritize the
        scheduler and trigger pump latency would suffer.
    -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
