# dist/launchd/com.marshall.nestty.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marshall.nestty.daemon</string>

    <!--
        ProgramArguments points at ~/.cargo/bin/nesttyd to match what
        scripts/install-macos.sh installs via `cargo install`. If you
        ran a system install that lands the binary elsewhere, edit
        the path here (or send the install script a PR with a
        $PROGRAM template).

        launchd does not expand ~ — that's why install-macos.sh
        rewrites the path with the user's HOME at copy time.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>HOME_PLACEHOLDER/.cargo/bin/nesttyd</string>
    </array>

    <!--
        RunAtLoad fires at launchd bootstrap (i.e. login). KeepAlive
        with no submap means "always running" — crash, manual kill,
        whatever; launchd restarts. Matches the systemd Restart=
        on-failure semantic close enough for parity.
    -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!--
        Pipe stdout + stderr to the user log dir so `tail
        ~/Library/Logs/nestty-daemon.{out,err}.log` is the macOS
        equivalent of journalctl user-level lookup for nestty-daemon.
    -->
    <key>StandardOutPath</key>
    <string>HOME_PLACEHOLDER/Library/Logs/nestty-daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_PLACEHOLDER/Library/Logs/nestty-daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <!--
        ProcessType=Interactive keeps the daemon at the same nice
        level as Aqua apps. Background would deprioritize the
        scheduler and trigger pump latency would suffer.
    -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
